# deploy/corterm-worker.service
[Unit]
Description=Corterm Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Uncomment for system-level service (not needed for user-level systemd)
# User={{USER}}
# Group={{GROUP}}
ExecStart={{INSTALL_DIR}}/corterm
WorkingDirectory={{HOME}}
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=corterm

# Environment
Environment=DOTNET_SYSTEM_NET_DISABLEIPV6=1

[Install]
WantedBy=multi-user.target
